# supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/fd/1
logfile_maxbytes=0
user=root

[program:x11]
priority=0
command=/usr/bin/Xtigervnc -desktop "noVNC" -localhost -rfbport 5900 -SecurityTypes None -AlwaysShared -AcceptKeyEvents -AcceptPointerEvents -AcceptSetDesktopSize -SendCutText -AcceptCutText :0
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:easy-novnc]
priority=0
command=/usr/local/bin/easy-novnc --addr :8080 --host localhost --port 5900 --no-url-password --novnc-params "resize=remote"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:openbox]
priority=1
command=/usr/bin/openbox-session
environment=DISPLAY=:0
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:tint2]
priority=1
command=tint2
environment=DISPLAY=:0
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:chrome]
priority=1
command=chrome freebitco.in --user-data-dir=/config/1 --start-maximized --disable-gpu --disable-software-rasterizer --no-sandbox --allow-running-insecure-content --disable-add-to-shelf --disable-background-networking --disable-background-timer-throttling --disable-backgrounding-occluded-windows --disable-breakpad --disable-checker-imaging --disable-client-side-phishing-detection --disable-component-extensions-with-background-pages --disable-datasaver-prompt --disable-desktop-notifications --disable-dev-shm-usage --disable-domain-reliability --disable-hang-monitor --disable-infobars --disable-ipc-flooding-protection --disable-notifications --disable-popup-blocking --disable-prompt-on-repost --disable-renderer-backgrounding --disable-setuid-sandbox --disable-site-isolation-trials --disable-sync --disable-web-security --enable-automation --ignore-certificate-errors --disable-logging --metrics-recording-only --mute-audio --no-default-browser-check --no-first-run --test-type --use-mock-keychain
environment=DISPLAY=:0
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

